# setup/systemd_units/home-xsup-magic_portal.automount
[Unit]
After=network-online.target
Before=remote-fs.target
[Automount]
Where=/home/xsup/magic_portal
TimeoutIdleSec=600
[Install]
WantedBy=multi-user.target
