# setup/systemd_units/renew_certificates.timer
[Unit]
Description=Renew certificates

[Timer]
Persistent=true
OnCalendar=*-*-* 1:23
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
